# buildroot/board/rockchip/rk3128/fs-overlay/etc/udev/rules.d/60-drm.rules
SUBSYSTEM=="drm", ACTION=="add|change", RUN+="/usr/bin/drm-hotplug.sh"
